# profile.d/bundler.sh
# Store bundled gems within the projects we download them for.
export BUNDLE_PATH=vendor/bundle
