# profile.d/bundler.sh
#!/bin/bash

export BUNDLE_PATH=vendor/bundle
