# profile.d/chruby.sh
# Automatically change Ruby versions when changing directories.
if [[ -d /usr/local/share/chruby ]]; then
    source /usr/local/share/chruby/chruby.sh
    source /usr/local/share/chruby/auto.sh
fi
