# bash/profile.d/qmk.sh
#!/bin/bash

export QMK_HOME=$HOME/config_files/keyboard/qmk
