# bash/profile.d/verboten_commands.sh
#!/bin/bash

# Prevent myself from using some commands that I don't want to use anymore.

alias nslookup='echo "Use dig!"'
alias traceroute='echo "Use mtr!"'
alias nano='echo "Use vim!"'
alias pico='echo "Use vim!"'
alias ack='echo "Use rg!"'
alias ag='echo "Use rg!"'
